# lcs-security/conf/cron.d/lcs-save-ldap
30 18 * * * root /usr/share/lcs/scripts/lcs-save-ldap.sh > /dev/null
