# lcs-bcdiweb/debian/postinst
#!/bin/sh
# postinst script for lcs-bcdiweb
#
# see: dh_installdeb(1)

set -e
preversion=$2
MODULE="bcdiweb"
VER=$(dpkg-query -f '${Version}' -W lcs-$MODULE)
INDICEMAJNBR="0"

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    if [ "$2" = "" ] || [ "$2" = "$VER" ]; then
    	 ### Install initiale
	mkdir -p /home/BcdiPlug/Data
	chown -R www-data /home/BcdiPlug
	mysql -e "INSERT INTO lcs_db.applis VALUES('', '$MODULE', '1', 'Plugin BCDI Web pour Lcs', '$INDICEMAJNBR', '$VER', 'bcdiweb/progweb', 'P');"
	 if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d apache2 restart
        else
                /etc/init.d/apache2 restart
        fi
      else
          ### mise a jour de lcs-bcdiweb
          if dpkg --compare-versions $preversion lt 2.4 ; then
           mysql -e "DELETE FROM lcs_db.params WHERE params.name = 'Bcdiweb' LIMIT 1;"
           mysql -e "INSERT INTO lcs_db.applis VALUES('', '$MODULE', '1', 'Plugin BCDI Web pour Lcs', '$INDICEMAJNBR', '$VER', 'bcdiweb', 'P');"		
	   
	  fi
	fi
	/usr/bin/mysql -e "UPDATE lcs_db.applis SET num_maj='$INDICEMAJNBR', version='$VER' WHERE name='$MODULE';"
	chown www-data -R /usr/share/lcs/Plugins/bcdiweb	
	;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


