# lcs-web/cron.d/lcsldapsav
30 18 * * * root /usr/share/lcs/scripts/lcs_sauve_ldap.sh > /dev/null
