# lcse3-monitor/debian/postinst
#!/bin/sh
# postinst script for gestEtab-monitor
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    install|configure)

# on rend executable les scripts
chmod +x /usr/sbin/gestEtabMonitor

# patch le cron
CRONNUMBER=$[ ( $RANDOM % 60 )  + 1 ];
/bin/sed -i'' -e "s#{CRONRAND}#${CRONNUMBER}#g" \
            "/etc/cron.d/lcse3-monitor"

if [ -e /usr/sbin/lcs-make-sudoers ]; then
	/usr/sbin/lcs-make-sudoers
else	
	##### ADD good line in /etc/sudoers.conf if necessary and restart sudo
	TEMOIN_CHGT_SUDO="n"
	# TEST=$(grep "NOPASSWD:GESTETABMONITOR" /etc/sudoers)
	if [ ! "$(grep "NOPASSWD:GESTETABMONITOR" /etc/sudoers)" ]; then
        	# Insertion de l'autorisation de lancement par www-data
        	sed -i 's/^\(www-data.*\)$/\1,NOPASSWD:GESTETABMONITOR/' /etc/sudoers
        	TEMOIN_CHGT_SUDO="y"
	fi

	# TEST=$(cat /etc/sudoers | grep "Cmnd_Alias GESTETABMONITOR")
	if [ ! "$(grep "Cmnd_Alias GESTETABMONITOR" /etc/sudoers)" ]; then
        	# On insère toute la ligne des commandes
        	sed -i 's|^\(# Cmnd alias specification.*\)$|\1\nCmnd_Alias GESTETABMONITOR=/usr/sbin/gestEtabMonitor|' /etc/sudoers
        	TEMOIN_CHGT_SUDO="y"
	fi

	if [ "$TEMOIN_CHGT_SUDO" = "y" ]; then
        	/etc/init.d/sudo restart
	fi
fi

if [ -e /etc/init.d/lcse3-monitor ] ; then 
    update-rc.d lcse3-monitor defaults 90 90
fi

#patch le fichier de conf
if [ -e /etc/gestEtab/etab.conf.patch ] ; then 
	patch /etc/gestEtab/etab.conf /etc/gestEtab/etab.conf.patch
	rm /etc/gestEtab/etab.conf.patch
fi

echo "lcse3-monitor est correctement installe. Si c'est votre premiere installation pensez à editer le fichier /etc/gestEtab/etab.conf pour configurer le systeme."
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
