# lcs-openvpn/debian/postinst
#!/bin/sh
# postinst script for lcs-openvpn
#
# see: dh_installdeb(1)

set -e

preversion=$2
MODULE="openvpn"
VER=$(dpkg-query -f '${Version}' -W lcs-$MODULE)
INDICEMAJNBR="0"

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
		if [ "$2" = "" ] || [ "$2" = "$VER" ]; then
			### Install initiale	
			/usr/share/lcs/sbin/lcs-openvpn-initial-config
			/usr/sbin/lcs-make-sudoers
			if [ -e /etc/default/openvpn ]; then
				sed -i  -e 's/^#\(AUTOSTART="all"\)/\1/' /etc/default/openvpn
			fi
		
			if [ -x /usr/sbin/invoke-rc.d ]; then
				invoke-rc.d openvpn restart
			else
				/etc/init.d/openvpn restart
			fi	
			
			#
			# Inscription du module dans applis.lcs_db
			#
			/usr/bin/mysql -e "INSERT INTO lcs_db.applis (id, name, value, descr, num_maj, version, chemin, type) VALUES ('', '$MODULE', '1', 'VPN par openvpn en SSL', '$INDICEMAJNBR', '$VER', '', 'S');"
		else
			### mise a jour de lcs-openvpn
			if [ -e /etc/default/openvpn ]; then
				sed -i  -e 's/^#\(AUTOSTART="all"\)/\1/' /etc/default/openvpn
        fi
			if dpkg --compare-versions $preversion lt 2.4 ; then
				# update mysql maj 2.4
				mysql -e "INSERT INTO lcs_db.params (id, name, value, srv_id, descr, cat) VALUES ('', 'vpnmasquerade', '0', 0, 'ip du LCS est utilise en NAT pour les clients de vpn', 11);"
				mysql -e "INSERT INTO lcs_db.params (id, name, value, srv_id, descr, cat) VALUES ('', 'vpnwins', 'aucun', 0, 'ip du serveur WINS generalement IP du SE3', 11);"
				mysql -e "INSERT INTO lcs_db.params (id, name, value, srv_id, descr, cat) VALUES ('', 'vpndnsclient', '0', 0, 'Utiliser le LCS en serveur DNS sur le VPN', 11);"
			fi

			echo "Mise ajour du paquet module LCS $MODULE de $2 vers $VER"
			# mise a jour dans lcs_db.applis
	    	TEST=`mysql -se "SELECT type FROM lcs_db.applis WHERE name='$MODULE';"`
	    	if 	[ "$TEST" != "S" ]; then
				/usr/bin/mysql -e "INSERT INTO lcs_db.applis (id, name, value, descr, num_maj, version, chemin, type) VALUES ('', '$MODULE', '1', 'VPN par openvpn en SSL', '$INDICEMAJNBR', '$VER', '', 'S');"
			else
				/usr/bin/mysql -e "UPDATE lcs_db.applis SET num_maj='$INDICEMAJNBR', version='$VER' WHERE name='$MODULE';"
	    	fi
	    	
	    	if dpkg --compare-versions $preversion lt 2.4.5 ; then
	    		/usr/sbin/lcs-make-sudoers
	    	fi
	    	
	    	# Activation de l'ip forward.
	    	echo "1" /proc/sys/net/ipv4/ip_forward
	    	
	    	# Si Debian Wheezy
	    	if  ! egrep -q "^6.0" /etc/debian_version;  then
	    		if [ ! -d /usr/share/php/smarty ]; then
	    			mkdir -p /usr/share/php/smarty
	    			ln -s /usr/share/php/smarty3 /usr/share/php/smarty/libs
	    	fi
	      # si pas de fichier config, recreation 
	      if [ ! -e /usr/share/lcs/openvpn/config_openvpn.inc.php ]; then
	      	HOSTNAME=`hostname`
				FULLDOMAINE=`hostname -f`
	      	cp /usr/share/lcs/openvpn/config.inc /usr/share/lcs/openvpn/config_openvpn.inc.php
				sed -i'' "s/@@HOSTNAME@@/$HOSTNAME/g" /usr/share/lcs/openvpn/config_openvpn.inc.php
				sed -i'' "s/@@FULLDOMAINE@@/$FULLDOMAINE/g" /usr/share/lcs/openvpn/config_openvpn.inc.php 
			fi	
			
	    	if dpkg --compare-versions $preversion lt 2.4.9.3 ; then
				/etc/init.d/apache2 restart  		
	    	fi			
			
			#
			# Fin de la mise a jour
			#
			echo "Termine."
    fi
    ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
	;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


