# lcs-security/debian/postinst
#!/bin/sh
# postinst script for lcs-security
#
# see: dh_installdeb(1)

set -e

preversion=$2
MODULE="security"
VER=$(dpkg-query -f '${Version}' -W lcs-$MODULE)
INDICEMAJNBR="0"

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    	if [ `grep -c rssh /etc/passwd` = 0 ]; then
    	### installation du compte rssh
    	useradd -K UID_MIN=100 -K UID_MAX=999 -K GID_MIN=100 -K GID_MAX=999 rssh
    	chown -R rssh\:rssh /home/rssh
    	chmod -R 700 /home/rssh/.ssh
    	cat >> /etc/inittab <<EOP
# Action on special keypress SSH inverse de secours XenEdu (ALT-UpArrow).
kb::kbrequest:/sbin/lcs-secours-ssh
EOP
		init q
	fi
	# Create check.log
	if [ ! -e /var/log/lcs/check.log ]; then
		if [ ! -d /var/log/lcs ]; then
			mkdir -p /var/log/lcs
			chown www-data:www-data /var/log/lcs
			chmod 750 /var/log/lcs 
		fi
		touch /var/log/lcs/check.log
		chmod 640 /var/log/lcs/check.log
		chown www-data:www-data /var/log/lcs/check.log
	fi
   # Restart fail2ban
	if [ -x /usr/sbin/invoke-rc.d ]; then
		invoke-rc.d fail2ban restart
	else
		/etc/init.d/fail2ban restart
	fi
	# Update sudoers lt 2.4.9.1
	if [ -n "$2" ] && dpkg --compare-versions "$2" lt 2.4.9.1; then	
				echo "update lcs-security and disable mod_apache_security"
		    	# Activation des logs sudo
		    	/usr/sbin/lcs-make-sudoers
		    	# Desactivation du mod security d'apache
		    	echo "le mod security d'apache est disponible."
		    	/usr/sbin/lcs-disable-apache-security 
	fi


    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


