# lcs-security/debian/postinst
#!/bin/sh
# postinst script for lcs-security
#
# see: dh_installdeb(1)

set -e

preversion=$2
MODULE="security"
VER=$(dpkg-query -f '${Version}' -W lcs-$MODULE)
INDICEMAJNBR="0"

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
	if [ -x /usr/sbin/invoke-rc.d ]; then
		invoke-rc.d fail2ban restart
	else
		/etc/init.d/fail2ban restart
	fi
	if dpkg --compare-versions $preversion lt 2.4.9 ; then
		    	# Activation des logs sudo
		    	/usr/sbin/lcs-make-sudoers
		    	# Activation du mod security d'apache
		    	/usr/sbin/lcs-enable-apache-security 
	fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


