# lcs-security/debian/postinst
#!/bin/sh
# postinst script for lcs-security
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
	# installation du compte rssh
	if ! getent passwd rssh >/dev/null ; then
		useradd --system rssh
		cp -a /usr/share/lcs-security/rssh /home/
		chown -R rssh:rssh /home/rssh
		chmod -R u=rwX,g=,o= /home/rssh/.ssh
		cat >> /etc/inittab <<EOP
# Action on special keypress SSH inverse de secours XenEdu (ALT-UpArrow).
kb::kbrequest:/sbin/lcs-secours-ssh
EOP
		telinit q
	fi
	# Installation de /root/.bashrc à l'installation initiale
	if [ -z "$2" ]; then
		cp -a /usr/share/lcs-security/.bashrc /root/.bashrc
	fi
	# Create check.log
	if [ ! -e /var/log/lcs/check.log ]; then
		if [ ! -d /var/log/lcs ]; then
			mkdir -p /var/log/lcs
			chown www-data:www-data /var/log/lcs
			chmod 750 /var/log/lcs
		fi
		touch /var/log/lcs/check.log
		chmod 640 /var/log/lcs/check.log
		chown www-data:www-data /var/log/lcs/check.log
	fi
	# Restart fail2ban
	invoke-rc.d fail2ban restart || true
	# Update sudoers lt 2.4.9.1
	if dpkg --compare-versions "$2" lt-nl 2.4.9.1; then	
		echo "update lcs-security and disable mod_apache_security"
		# Activation des logs sudo
		lcs-make-sudoers
		# Desactivation du mod security d'apache
		echo "le mod security d'apache est disponible."
		lcs-disable-apache-security
	fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


