# lcs-bcdiweb/conf/cron.d/bcdiweb
0 13 * * 0-5 www-data  /usr/share/lcs/bcdiweb/scripts/syncbcdi.sh
